# infra/systemd/liveops-backup.service
[Unit]
Description=Live Ops PostgreSQL verified backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/live-ops/current/infra/scripts/backup_production.sh
